# src/main/assembly/dist/install/dd-manage-deposit.service
[Unit]
Description=Dd Manage Deposit Service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-manage-deposit/bin/dd-manage-deposit server /etc/opt/dans.knaw.nl/dd-manage-deposit/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143

User=dd-manage-deposit
Group=dd-manage-deposit

[Install]
WantedBy=multi-user.target
